# deploy/ruuvi-prometheus.timer
[Unit]
Description=Start scraper after boot on a delay

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
